# ops/gotouhou-hourly-progress.service
[Unit]
Description=Send gotouhou hourly development progress email

[Service]
Type=oneshot
WorkingDirectory=/root/gotouhou/docs
EnvironmentFile=-/etc/gotouhou/progress-mail.env
ExecStart=/bin/sh /root/gotouhou/docs/ops/hourly_progress_runner.sh
