# ops/gotouhou-hourly-progress.service
[Unit]
Description=Send gotouhou hourly development progress email

[Service]
Type=oneshot
WorkingDirectory=/root/gotouhou/docs
EnvironmentFile=-/etc/gotouhou/progress-mail.env
ExecStartPre=-/usr/bin/python3 /root/gotouhou/docs/ops/agent_watchdog.py
ExecStart=/usr/bin/python3 /root/gotouhou/docs/ops/hourly_progress_mail.py --brief
