# ops/gotouhou-hourly-progress.service
[Unit]
Description=Send gotouhou three-hour development progress email

[Service]
Type=oneshot
WorkingDirectory=/root/gotouhou/docs
EnvironmentFile=-/etc/gotouhou/progress-mail.env
Environment=HOME=/root
Environment=XDG_CONFIG_HOME=/root/.config
Environment=GOCACHE=/root/.cache/go-build
Environment=GOPATH=/root/go
TimeoutStartSec=20min
ExecStart=/bin/sh /root/gotouhou/docs/ops/hourly_progress_runner.sh
