# ops/gotouhou-hourly-progress.timer
[Unit]
Description=Run gotouhou three-hour development progress email

[Timer]
OnCalendar=*-*-* 00/3:00:00
Persistent=true
Unit=gotouhou-hourly-progress.service

[Install]
WantedBy=timers.target
