# ops/gotouhou-hourly-progress.timer
[Unit]
Description=Run gotouhou hourly development progress email

[Timer]
OnCalendar=hourly
Persistent=true
Unit=gotouhou-hourly-progress.service

[Install]
WantedBy=timers.target
